# scotia/deploy/scotiad.service
[Unit]
Description=Scotia daemon — session registry and notification hub
After=network.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/scotiad
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
